# src/init.d/keusbd
#!/sbin/runscript
#
# init.d/keusbd -- cценарий управления демоном keusbd
#
# Разработал: Муковников М. Ю.
# Версия:     1.0.0
# Изменён:    16.04.2010
#

depend() {
	need net
}

start() {
	ebegin "Starting keusbd"
	start-stop-daemon --start --quiet --exec /usr/bin/keusbd -- --daemon
	eend $?
}

stop() {
	ebegin "Stopping keusbd"
	start-stop-daemon --stop --quiet --pidfile /var/run/keusbd.pid
	eend $?
}

